# files/systemd/cephcollectorImages.service
[Unit]
Description=Ceph collector - pool images

[Service]
ExecStart=/usr/local/sbin/ceph-collector.sh images

[Install]
WantedBy=multi-user.target
