# templates/cephcollector.service.epp
<%- | String  $pool_name,
| -%>
[Unit]
Description=Ceph collector - pool <%= $pool_name %>

[Service]
ExecStart=/usr/local/sbin/ceph-collector.sh <%= $pool_name %> 

[Install]
WantedBy=multi-user.target
