# templates/systemd.service.epp
<%- | String $description,
      String $command
| -%>
[Unit]
Description=<%= $description %>

[Service]
ExecStart=<%= $command %>

[Install]
WantedBy=multi-user.target
